# supervisor.conf
[program:vncserver]
command=vncserver -geometry 1280x800 :1
user=ubuntu
